# call-home.service
[Unit]
Description=Run script to automatically backup home configs

[Service]
Type=simple
ExecStart=/home/vermeer/call-home.sh
